# installer/debian/oco-agent/DEBIAN/prerm
#!/bin/bash

if [ "$1" != "upgrade" ]; then
    systemctl stop oco-agent
    systemctl disable oco-agent
fi
