# installer/debian/oco-agent/DEBIAN/prerm
#!/bin/bash

systemctl stop oco-agent
systemctl disable oco-agent

